# addons/etc/update-motd.d/00-figlet
#!/bin/sh

. /etc/update-motd.d/colors

printf "\n"$LIGHT_RED
figlet "  "RuyiSDK
printf "                   -- Presented by ISCAS"
printf "\n"
printf $NONE
printf "\n"
